# shared/osquery/osquery.conf
{
  "options": {
    "config_plugin": "filesystem",
    "logger_plugin": "filesystem",
    "logger_path": "/var/log/osquery",
    "disable_logging": "false",
    "log_result_events": "true",
    "schedule_splay_percent": "10",
    "pidfile": "/var/osquery/osquery.pidfile",
    "events_expiry": "3600",
    "database_path": "/var/osquery/osquery.db",
    "verbose": "false",
    "worker_threads": "2",
    "enable_monitor": "true",
    "disable_events": "false",
    "disable_audit": "false",
    "audit_allow_config": "true",
    "audit_allow_sockets": "true",
    "host_identifier": "hostname",
    "enable_syslog": "true",
    "audit_allow_process_events": "true"
  },

  "schedule": {
    "process_events": {
      "query": "SELECT pid, path, cmdline, uid, gid, parent, cwd, start_time, CASE WHEN path LIKE '/tmp/%' OR path LIKE '/dev/shm/%' OR cmdline LIKE '%wget%' OR cmdline LIKE '%curl%' OR cmdline LIKE '%nc %' OR cmdline LIKE '%python%reverse%' THEN 'SUSPICIOUS' ELSE 'normal' END AS risk_level FROM processes WHERE start_time > (strftime('%s','now') - 60);",
      "interval": 15,
      "description": "Monitoreo de procesos cada 15 segundos - detecta ejecuciones sospechosas"
    },
    "listening_ports": {
      "query": "SELECT DISTINCT p.name, p.path, lp.port, lp.protocol, lp.address, u.username FROM listening_ports lp JOIN processes p ON lp.pid = p.pid JOIN users u ON p.uid = u.uid WHERE lp.port != 0;",
      "interval": 30,
      "description": "Puertos en escucha - detecta backdoors y shells reversas"
    },
    "network_connections": {
      "query": "SELECT DISTINCT p.name, p.path, p.cmdline, ps.remote_address, ps.remote_port, ps.local_port, ps.state, u.username FROM process_open_sockets ps JOIN processes p ON ps.pid = p.pid JOIN users u ON p.uid = u.uid WHERE ps.remote_address != '' AND ps.remote_address != '127.0.0.1' AND ps.remote_address != '::1' AND ps.state = 'ESTABLISHED';",
      "interval": 20,
      "description": "Conexiones de red activas - detecta comunicaciones C2"
    },
    "file_events_tmp": {
      "query": "SELECT target_path, action, md5, sha256, uid, gid, atime, mtime, ctime FROM file_events WHERE target_path LIKE '/tmp/%' OR target_path LIKE '/dev/shm/%' OR target_path LIKE '/var/tmp/%';",
      "interval": 10,
      "description": "Eventos de archivos en directorios temporales - detecta payloads"
    },
    "crontab_monitor": {
      "query": "SELECT * FROM crontab;",
      "interval": 60,
      "description": "Monitoreo de tareas programadas - detecta persistencia"
    },
    "shell_history": {
      "query": "SELECT uid, username, command, history_file FROM shell_history WHERE command LIKE '%wget%' OR command LIKE '%curl%' OR command LIKE '%chmod +x%' OR command LIKE '%python%' OR command LIKE '%nc %' OR command LIKE '%base64%' OR command LIKE '%eval%';",
      "interval": 30,
      "description": "Historial de shell sospechoso - detecta actividad maliciosa"
    },
    "user_ssh_keys": {
      "query": "SELECT uid, path, encrypted FROM user_ssh_keys;",
      "interval": 120,
      "description": "Monitoreo de claves SSH - detecta movimiento lateral"
    },
    "logged_in_users": {
      "query": "SELECT type, user, host, time, pid FROM logged_in_users;",
      "interval": 30,
      "description": "Usuarios conectados - detecta accesos no autorizados"
    },
    "suspicious_modules": {
      "query": "SELECT name, path, size FROM kernel_modules WHERE name NOT IN ('ext4','xfs','nfs','overlay','bridge','br_netfilter','ip_tables','iptable_nat','nf_nat');",
      "interval": 120,
      "description": "Módulos del kernel - detecta rootkits"
    },
    "dns_resolvers": {
      "query": "SELECT * FROM dns_resolvers;",
      "interval": 300,
      "description": "Resolvers DNS - detecta DNS hijacking"
    },
    "suid_binaries": {
      "query": "SELECT path, username, permissions, mtime FROM suid_bin;",
      "interval": 300,
      "description": "Binarios SUID - detecta escalación de privilegios"
    },
    "docker_detection": {
      "query": "SELECT pid, name, path, cmdline FROM processes WHERE name IN ('docker','dockerd','containerd','runc');",
      "interval": 60,
      "description": "Detección de contenedores - detecta container escape"
    },
    "yara_scan_tmp": {
      "query": "SELECT path, count, matches, strings, tags FROM yara WHERE path LIKE '/tmp/%' AND sigfile = '/etc/osquery/yara/malware.yar';",
      "interval": 60,
      "description": "Escaneo YARA en /tmp - detecta malware conocido"
    },
    "process_open_files": {
      "query": "SELECT p.pid, p.name, p.path, pof.path AS open_file FROM processes p JOIN process_open_files pof ON p.pid = pof.pid WHERE pof.path LIKE '/tmp/%' OR pof.path LIKE '/dev/shm/%' OR pof.path LIKE '%.sh' OR pof.path LIKE '%.py';",
      "interval": 30,
      "description": "Archivos abiertos por procesos - detecta payloads en ejecución"
    },
    "authorized_keys_monitor": {
      "query": "SELECT * FROM authorized_keys;",
      "interval": 60,
      "description": "Monitoreo de authorized_keys - detecta persistencia SSH"
    },
    "arp_cache": {
      "query": "SELECT address, mac, interface FROM arp_cache;",
      "interval": 60,
      "description": "Cache ARP - detecta ARP spoofing"
    }
  },

  "file_paths": {
    "tmp_files": ["/tmp/%%", "/var/tmp/%%", "/dev/shm/%%"],
    "binaries": ["/usr/bin/%%", "/usr/sbin/%%", "/usr/local/bin/%%"],
    "ssh_configs": ["/root/.ssh/%%", "/home/%/.ssh/%%"],
    "cron_files": ["/etc/cron.d/%%", "/var/spool/cron/%%", "/etc/crontab"],
    "startup_files": ["/etc/init.d/%%", "/etc/systemd/system/%%"],
    "web_shells": ["/var/www/%%", "/srv/%%"]
  },

  "events": {
    "disable_subscribers": ["user_events"]
  },

  "packs": {
    "threat-hunting": "/etc/osquery/packs/threat-hunting.conf"
  }
}
