# shared/osquery/packs/threat-hunting.conf
{
  "queries": {
    "reverse_shell_detection": {
      "query": "SELECT p.pid, p.name, p.path, p.cmdline, p.parent, pp.name AS parent_name, pp.cmdline AS parent_cmdline FROM processes p LEFT JOIN processes pp ON p.parent = pp.pid WHERE (p.cmdline LIKE '%/dev/tcp%' OR p.cmdline LIKE '%bash -i%' OR p.cmdline LIKE '%python%import socket%' OR p.cmdline LIKE '%perl%socket%' OR p.cmdline LIKE '%nc -e%' OR p.cmdline LIKE '%ncat%');",
      "interval": 10,
      "description": "Detecta reverse shells activas - T1059.004",
      "value": "Detección de shells reversas en tiempo real"
    },
    "crypto_miner_detection": {
      "query": "SELECT pid, name, path, cmdline FROM processes WHERE cmdline LIKE '%stratum%' OR cmdline LIKE '%xmrig%' OR cmdline LIKE '%minerd%' OR cmdline LIKE '%cpuminer%' OR name LIKE '%miner%';",
      "interval": 30,
      "description": "Detecta cryptominers - T1496",
      "value": "Detección de minería de criptomonedas"
    },
    "data_exfiltration": {
      "query": "SELECT p.pid, p.name, p.cmdline, ps.remote_address, ps.remote_port FROM processes p JOIN process_open_sockets ps ON p.pid = ps.pid WHERE (p.cmdline LIKE '%tar%' OR p.cmdline LIKE '%zip%' OR p.cmdline LIKE '%base64%') AND ps.remote_address != '' AND ps.remote_address != '127.0.0.1';",
      "interval": 15,
      "description": "Detecta exfiltración de datos - T1048",
      "value": "Detección de transferencia de datos sospechosa"
    },
    "persistence_mechanisms": {
      "query": "SELECT path, filename, mtime FROM file WHERE (path LIKE '/etc/cron.d/%' OR path LIKE '/etc/init.d/%' OR path LIKE '/etc/systemd/system/%' OR path LIKE '/home/%/.bashrc' OR path LIKE '/root/.bashrc') AND mtime > (strftime('%s','now') - 300);",
      "interval": 30,
      "description": "Detecta mecanismos de persistencia recientes - T1053/T1546",
      "value": "Detección de persistencia"
    },
    "privilege_escalation_suid": {
      "query": "SELECT path, username, groupname, permissions, mtime FROM suid_bin WHERE mtime > (strftime('%s','now') - 3600);",
      "interval": 60,
      "description": "Detecta nuevos binarios SUID - T1548.001",
      "value": "Detección de escalación de privilegios via SUID"
    },
    "encoded_commands": {
      "query": "SELECT pid, name, path, cmdline FROM processes WHERE cmdline LIKE '%base64%' OR cmdline LIKE '%eval%' OR cmdline LIKE '%-enc %' OR cmdline LIKE '%frombase64%' OR cmdline LIKE '%decode%';",
      "interval": 15,
      "description": "Detecta comandos codificados - T1027",
      "value": "Detección de ofuscación de comandos"
    },
    "lateral_movement_ssh": {
      "query": "SELECT p.pid, p.name, p.cmdline, p.uid, u.username, ps.remote_address, ps.remote_port FROM processes p JOIN process_open_sockets ps ON p.pid = ps.pid JOIN users u ON p.uid = u.uid WHERE p.name = 'ssh' AND ps.remote_port = '22' AND ps.state = 'ESTABLISHED';",
      "interval": 15,
      "description": "Detecta movimiento lateral via SSH - T1021.004",
      "value": "Detección de movimiento lateral"
    },
    "webshell_detection": {
      "query": "SELECT path, filename, size, mtime FROM file WHERE (path LIKE '/var/www/%.php' OR path LIKE '/var/www/%.jsp' OR path LIKE '/var/www/%.asp') AND mtime > (strftime('%s','now') - 3600);",
      "interval": 60,
      "description": "Detecta web shells recientes - T1505.003",
      "value": "Detección de web shells"
    }
  }
}
