# endpoints/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true

[program:rsyslog]
command=/usr/sbin/rsyslogd -n
autostart=true
autorestart=true

[program:osqueryd]
command=/usr/bin/osqueryd --config_path=/etc/osquery/osquery.conf --logger_path=/var/log/osquery --pidfile=/var/osquery/osquery.pidfile --database_path=/var/osquery/osquery.db --disable_events=false --disable_audit=false --host_identifier=hostname
autostart=true
autorestart=true
startretries=5
startsecs=5

[program:filebeat]
command=/usr/bin/filebeat -e -c /etc/filebeat/filebeat.yml
autostart=true
autorestart=true
startretries=5

[program:vnc]
command=bash -c "vncserver :1 -geometry 1280x800 -depth 24 -localhost no -fg"
autostart=true
autorestart=true

[program:novnc]
command=/opt/noVNC/utils/websockify/run --web /opt/noVNC 6901 localhost:5901
autostart=true
autorestart=true
startretries=5

[program:process_monitor]
command=/opt/monitoring/process_monitor.sh
autostart=true
autorestart=true

[program:network_monitor]
command=/opt/monitoring/network_monitor.sh
autostart=true
autorestart=true

[program:file_monitor]
command=/opt/monitoring/file_monitor.sh
autostart=true
autorestart=true
